# infra/ib-gateway.service
[Unit]
Description=IB Gateway (headless via IBC + Xvfb)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
Environment="DISPLAY=:99"
Environment="LD_LIBRARY_PATH=/usr/lib/jvm/java-17-openjdk-amd64/lib/server:/usr/lib/jvm/java-17-openjdk-amd64/lib"
ExecStart=/opt/ibc/launchgateway.sh
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
