# infra/ib-gateway.service
[Unit]
Description=IB Gateway (headless via IBC + Xvfb)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
Environment="DISPLAY=:99"
ExecStart=/opt/ibc/launchgateway.sh
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
